# systemd/atd.service
[Unit]
Description=AT Daemon

[Service]
Type=forking
PIDFile=/run/atd.pid
ExecStart=/usr/sbin/atd

[Install]
WantedBy=multi-user.target
