# systemd/getty@.service
.include /usr/lib/systemd/system/getty@.service

[Service]
TTYVTDisallocate=no
