# systemd/hdparm-mobile.service
[Unit]
Description=Apply hdparm settings

[Service]
Type=simple
ExecStart=/home/void/bin/hdparm-mobile
TimeoutSec=0
