# systemd/mpd-as-user@.service
.include /usr/lib/systemd/system/mpd.service

[Unit]
Description=Music Player Daemon (as user %i)

[Service]
User=%i
