# systemd/netbook-powersave.service
[Unit]
Description=Apply powersave settings

[Service]
Type=simple
ExecStart=/home/void/bin/netbook-powersave
TimeoutSec=0
