# systemd/temperature-log.service
[Unit]
Description=Temperature Log

[Service]
Type=simple
ExecStart=/home/void/bin/temperature-log

[Install]
WantedBy=multi-user.target
